# docker/django/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/code/uwsgi.ini
stopsignal = QUIT

[program:nginx-app]
command = /usr/sbin/nginx
stopsignal = QUIT